# server/chroot_files/persistent.conf
description "Mount /persistent"

start on filesystem

script
  DEV=$(findfs LABEL=persistent)
  if test "$?" = "0"; then
    mkdir /persistent
    mount "${DEV}" /persistent
  fi
end script
